# macos/Library/LaunchAgents/com.evgenii.tailscale-serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.evgenii.tailscale-serve</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/evgenii/.local/bin/tailscale-serve</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>15</integer>
    <key>StandardOutPath</key>
    <string>/tmp/tailscale-serve.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tailscale-serve.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/evgenii/.local/bin</string>
    </dict>
</dict>
</plist>
