# macos/Library/LaunchAgents/com.evgenii.tmux.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.evgenii.tmux</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/opt/homebrew/bin/tmux start-server 2&gt;/dev/null; PID=$(/opt/homebrew/bin/tmux display-message -p &quot;#{pid}&quot;); while kill -0 &quot;$PID&quot; 2&gt;/dev/null; do sleep 5; done</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/tmp/tmux.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tmux.error.log</string>
</dict>
</plist>
