# macos/Library/LaunchAgents/com.evgenii.opencode.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.evgenii.opencode</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/evgenii/.opencode/bin/opencode</string>
        <string>serve</string>
        <string>--hostname</string>
        <string>0.0.0.0</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/opencode.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/opencode.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/evgenii/.opencode/bin:/Users/evgenii/.local/bin</string>
        <key>HOME</key>
        <string>/Users/evgenii</string>
    </dict>
</dict>
</plist>
