# opencode/.config/systemd/user/opencode-tailscale.service
[Unit]
Description=Tailscale Serve for OpenCode
Documentation=https://tailscale.com/kb/1242/tailscale-serve
After=opencode.service
BindsTo=opencode.service
PartOf=opencode.service

[Service]
Type=exec
# Wait for opencode to be ready before starting tailscale serve
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/tailscale serve http://127.0.0.1:4096
ExecStop=/usr/bin/tailscale serve off
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
